# install/systemd_files/time-lapse.service
[Unit]
Description=time-lapse

[Service]
Type=oneshot
ExecStart=/bin/bash take_picture.sh
WorkingDirectory=/home/pi/picam_wrapper/
User=pi
