# install/systemd_files/time-lapse.timer
[Unit]
Description=time-lapse picture

[Timer]
OnCalendar=00/1:00:00

[Install]
WantedBy=timers.target
