# install/systemd_files/time-lapse.timer
[Unit]
Description=time-lapse picture

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
